# edge/jetson/smart-classroom-edge.service
# systemd unit for the Jetson edge recognition service.
# Installed by install.sh (versioned release model): code lives in
# /opt/smart-classroom-edge/current (symlink -> releases/<version>) with its own
# venv; config in /opt/smart-classroom-edge/.env. Logs: journalctl -u smart-classroom-edge -f
# Rollback = repoint `current` at an older release + `systemctl restart`.

[Unit]
Description=Smart Classroom Jetson edge face recognition
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jetson
# `current` resolves at start, so a rolled-back symlink takes effect on restart.
WorkingDirectory=/opt/smart-classroom-edge/current
EnvironmentFile=/opt/smart-classroom-edge/.env
ExecStart=/opt/smart-classroom-edge/current/venv/bin/python /opt/smart-classroom-edge/current/recognize_service.py
Restart=always
RestartSec=5
# Nano power: keep it from being OOM-killed before others.
OOMScoreAdjust=-200

[Install]
WantedBy=multi-user.target
